# r-server/supervisord.conf
[supervisord]
nodaemon=true

[program:r_serve]
command=R CMD BATCH /run-rserve.R
autostart=true
redirect_stderr=true
stdout_logfile=/dev/stderr
autorestart=unexpected
exitcodes=0
